# crontab

# m h  dom mon dow   command

*/15 * * * * /usr/local/bin/php /opt/roardaemon/readrooms.php
